# etc/lava-dispatcher-host.socket
[Unit]
Description=lava-dispatcher-host socket

[Socket]
ListenStream=/run/lava-dispatcher-host.sock

[Install]
WantedBy=sockets.target
