# etc/lava-dispatcher-host.service
[Unit]
Description=lava-dispatcher-host daemon

[Service]
ExecStart=/usr/bin/lava-dispatcher-host-server
